# terraform_ha_stack/modules/instances/scripts/chef/cookbooks/configure-app/templates/supervisor.conf.erb
[program:flaskapp]
directory=/var/flask_app
command=/usr/local/bin/gunicorn -w 3 run:app --bind 0.0.0.0
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
